# crontab
* * * * * ~/convox_assist/idle_monitor.sh
